# sv/sxhkd/run
#!/usr/bin/execlineb -P
importas home HOME
cd $home
fdmove -c 2 1
sxhkd
