# sv/dwm/run
#!/usr/bin/execlineb -P
importas home HOME
cd $home
dwm
